# python/server.conf
description "Example JSON HTTP server"
author      "nitaku - matteo.abrate@gmail.com"

start on started mountall
stop on shutdown

respawn
respawn limit 99 5

script
  exec sudo -u www-data /usr/bin/python /home/angelo/Software/esp8266/project/Contatore_luce_nonos/python/angelo_server.py 5008 >> /dev/null 2>> /dev/null

end script

post-start script
end script
